# 案件/fusion_hrmos/hrmos_automation/com.tokumori.hrmos.weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.tokumori.hrmos.weekly</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>/Users/atsuyasato/Claude AI/hrmos_automation/weekly_run.py</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/atsuyasato/Claude AI</string>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key><integer>3</integer>
    <key>Hour</key><integer>10</integer>
    <key>Minute</key><integer>30</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/atsuyasato/Claude AI/hrmos_automation/weekly.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/atsuyasato/Claude AI/hrmos_automation/weekly.err.log</string>
</dict>
</plist>
